# ansible/collections/ansible_collections/homelab/apps/roles/nut/templates/nut-prometheus-exporter.service.j2
[Unit]
Description=NUT to Prometheus Exporter
After=nut-server.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/nut_to_prometheus.sh
