# ansible/collections/ansible_collections/homelab/apps/roles/nut/templates/nut-prometheus-exporter.timer.j2
[Unit]
Description=Run NUT to Prometheus Exporter every 15 seconds

[Timer]
OnBootSec=10
OnUnitActiveSec=15
AccuracySec=1s

[Install]
WantedBy=timers.target
